# auto_down.cron
*/1 * * * *  /usr/src/OneBuyAPI/auto_down.sh